# deploy/mywebapp.socket
[Unit]
Description=Task Tracker Socket Connection

[Socket]
ListenStream=5000
NoDelay=true

[Install]
WantedBy=sockets.target